# etc/init.d/wipe_data
#!/bin/sh /etc/rc.common

START=96
PROG=/usr/bin/wipe_data

start() {
    $PROG &
}
