# usr/lib/opkg/info/device_manager.list
/etc/init.d/refresh_device_status
/usr/bin/device_query
/etc/init.d/device_manager
/usr/bin/device_manager
/usr/bin/ubus_call
