# usr/lib/opkg/info/lvgl_klipper.list
/etc/sysConfig/UIResource/R2.8/key_del.png
/etc/sysConfig/defData/error_code_map.json
/etc/sysConfig/UIResource/R2.8/setting_guide.png
/usr/bin/upgrade-server
/etc/sysConfig/UIResource/R2.8/Chassis.png
/etc/sysConfig/UIResource/R2.8/fan_off.png
/etc/sysConfig/defData/SourceHanSansSC-Bold.otf
/etc/sysConfig/UIResource/R2.8/ZMove.png
/etc/sysConfig/UIResource/R2.8/test/Test15.png
/etc/sysConfig/UIResource/R2.8/setting_about.png
/etc/sysConfig/UIResource/R2.8/file_prev.png
/etc/sysConfig/UIResource/R2.8/prev_page_enable.png
/etc/sysConfig/UIResource/K1_Max/img_humidity.png
/etc/sysConfig/UIResource/R2.8/noedit_texture_100x51.png
/etc/sysConfig/UIResource/K1/zoffset_Add_Normal.png
/etc/sysConfig/UIResource/K1/selfchecking_1.png
/etc/sysConfig/UIResource/R2.8/prev_page_disable.png
/etc/sysConfig/UIResource/R2.8/file_back.png
/etc/sysConfig/UIResource/R2.8/test/Test24.png
/etc/sysConfig/UIResource/K1/img_cfs_humidity_normal.png
/etc/sysConfig/UIResource/R2.8/clock2.png
/etc/sysConfig/audio/chinese/1.button_beep.wav
/etc/sysConfig/UIResource/R2.8/history_hook.png
/etc/sysConfig/UIResource/R2.8/test/Test3.png
/etc/sysConfig/UIResource/K1/zoffset_Sub_Normal.png
/etc/sysConfig/UIResource/R2.8/setting_camera.png
/etc/sysConfig/UIResource/K1/NotVisible.png
/etc/sysConfig/UIResource/K1/SortDownUnSelect.png
/etc/sysConfig/UIResource/K1/excluceObjectIcon.png
/etc/sysConfig/UIResource/K1/errorTips.png
/etc/sysConfig/UIResource/R2.8/home_noz.png
/etc/sysConfig/UIResource/R2.8/home_ladia.png
/etc/sysConfig/UIResource/R2.8/img_f011_page_ctrl.png
/etc/sysConfig/UIResource/R2.8/img_f011_resume.png
/etc/sysConfig/UIResource/R2.8/EstateRunoutOut.png
/etc/sysConfig/defData/file_preview.png
/etc/sysConfig/defData/folder_icon.png
/etc/sysConfig/UIResource/R2.8/img_f011_print_speed.png
/etc/sysConfig/UIResource/K1/SortDownSelect.png
/etc/sysConfig/UIResource/R2.8/warningTips.png
/usr/bin/audio-server
/etc/sysConfig/UIResource/K1/localexporttips.png
/etc/sysConfig/UIResource/K1/img_humidity_arrow.png
/etc/sysConfig/UIResource/R2.8/setting_selftest.png
/etc/sysConfig/UIResource/R2.8/XYMove.png
/usr/bin/Monitor
/etc/init.d/app
/etc/sysConfig/UIResource/K1/img_boot0.png
/etc/sysConfig/audio/chinese/1.sound_button_smart.wav
/etc/sysConfig/UIResource/R2.8/test/Test16.png
/etc/sysConfig/UIResource/R2.8/CFS_manual.png
/etc/sysConfig/UIResource/R2.8/home_logo.png
/etc/sysConfig/UIResource/R2.8/home_FAQ.png
/etc/sysConfig/UIResource/K1/selfchecking_2.png
/etc/sysConfig/defData/Auto_pressure_advance_testpadvance_product_test.gcode
/etc/sysConfig/UIResource/R2.8/CFS_shelf_A.png
/etc/sysConfig/UIResource/R2.8/img_f011_page_setting.png
/etc/sysConfig/UIResource/K1/bed.png
/etc/sysConfig/UIResource/R2.8/img_f011_refresh_qrcode.png
/etc/sysConfig/UIResource/R2.8/setting_reset.png
/etc/sysConfig/UIResource/R2.8/test/Test4.png
/usr/bin/web-server
/etc/sysConfig/UIResource/K1_Max/img_setup.png
/etc/sysConfig/UIResource/R2.8/setting_screen_off.png
/etc/sysConfig/UIResource/R2.8/icon_speed_mode.png
/etc/sysConfig/UIResource/R2.8/setting_check_version.png
/etc/sysConfig/UIResource/K1/time_used_1.png
/etc/sysConfig/defData/laser_offset_correction_0.2mm.gcode
/etc/sysConfig/UIResource/R2.8/home_file.png
/etc/sysConfig/UIResource/R2.8/img_f011_title_tab_right.png
/etc/sysConfig/UIResource/R2.8/file_clock.png
/etc/sysConfig/UIResource/K1/usedtimeicon.png
/etc/sysConfig/UIResource/K1/TilteLeftstyle.png
/etc/sysConfig/UIResource/R2.8/net_wifi.png
/etc/sysConfig/UIResource/R2.8/test/Test17.png
/etc/sysConfig/UIResource/R2.8/kb_lowercase_icon.png
/etc/sysConfig/UIResource/K1/Subtract.png
/etc/sysConfig/UIResource/R2.8/kb_uppercase_icon.png
/etc/sysConfig/UIResource/K1/TitleRightstyle.png
/etc/sysConfig/UIResource/R2.8/error_report.png
/etc/sysConfig/UIResource/K1/objectStateFocus.png
/etc/sysConfig/UIResource/K1/selfchecking_3.png
/etc/sysConfig/UIResource/K1_Max/K1MAX_Machine.png
/etc/sysConfig/UIResource/R2.8/CFS_shelf_B.png
/etc/sysConfig/UIResource/R2.8/home_camera.png
/etc/sysConfig/UIResource/R2.8/detail.png
/etc/sysConfig/UIResource/K1_Max/img_nozzle.png
/etc/sysConfig/UIResource/K1_Max/img_eye.png
/etc/sysConfig/UIResource/K1/print_success.png
/etc/sysConfig/UIResource/K1/Visible.png
/etc/sysConfig/UIResource/R2.8/img_debris_remove_F009.png
/etc/sysConfig/UIResource/R2.8/img_f011_light_closed.png
/etc/sysConfig/UIResource/R2.8/test/Test5.png
/etc/sysConfig/UIResource/K1_Max/img_vector_1120.png
/etc/sysConfig/UIResource/K1/K1_CheckTips.png
/etc/sysConfig/UIResource/K1/fila_used_1.png
/etc/sysConfig/UIResource/R2.8/bar.png
/etc/sysConfig/defData/file_print_photo.png
/etc/sysConfig/UIResource/R2.8/img_debris_remove.png
/etc/sysConfig/defData/line_width_height.gcode
/etc/sysConfig/UIResource/K1/zoffset_Add_Disable.png
/etc/sysConfig/UIResource/R2.8/arrow_right_8x9.png
/etc/sysConfig/UIResource/R2.8/fan_on.png
/etc/sysConfig/UIResource/R2.8/checkbox_grey.png
/etc/sysConfig/UIResource/K1/layers.png
/etc/sysConfig/defData/material_database.json
/etc/sysConfig/UIResource/R2.8/error_export.png
/etc/sysConfig/audio/chinese/1.sound_button_crisp.wav
/etc/sysConfig/UIResource/R2.8/test/Test10.png
/etc/sysConfig/UIResource/R2.8/resume_stop.png
/etc/sysConfig/UIResource/R2.8/test/Test18.png
/usr/bin/burn-server
/etc/sysConfig/UIResource/K1/SortIcon.png
/etc/sysConfig/UIResource/R2.8/img_f011_humi.png
/etc/sysConfig/UIResource/R2.8/EstateRunoutIn.png
/etc/sysConfig/UIResource/K1/zoffset_Sub_Disable.png
/etc/sysConfig/UIResource/K1/objectStateCancle.png
/etc/sysConfig/UIResource/K1/objectStatePrinting.png
/etc/sysConfig/UIResource/R2.8/complete.png
/etc/sysConfig/UIResource/K1/nuzzle.png
/etc/sysConfig/UIResource/R2.8/img_f011_exclude_objects_flag.png
/etc/sysConfig/UIResource/K1_Max/img_polygon_25.png
/etc/sysConfig/UIResource/R2.8/Stepper.png
/etc/sysConfig/UIResource/R2.8/test/Test6.png
/etc/sysConfig/UIResource/R2.8/selftest_right.png
/etc/sysConfig/UIResource/K1/img_refill_dir.png
/etc/sysConfig/UIResource/K1/img_humidity_normal_rangle.png
/etc/sysConfig/UIResource/R2.8/next.png
/etc/sysConfig/UIResource/R2.8/pause_stop.png
/etc/sysConfig/defData/Auto_strain_gauge_pa_valid_print.gcode
/etc/sysConfig/UIResource/K1/cameraicon.png
/etc/sysConfig/UIResource/R2.8/CFS_illustration_A.png
/etc/sysConfig/UIResource/R2.8/home_chart_y.png
/etc/sysConfig/defData/file_icon.png
/etc/sysConfig/UIResource/R2.8/home_bed.png
/usr/bin/wifi-server
/etc/sysConfig/UIResource/R2.8/selftest_wrong.png
/etc/sysConfig/UIResource/R2.8/opt_mark.png
/etc/sysConfig/defData/Auto_strain_gauge_pa_config.gcode
/etc/sysConfig/UIResource/R2.8/wifi_state_1.png
/etc/sysConfig/UIResource/K1/lightclose.png
/etc/sysConfig/UIResource/K1/icon_chamber.png
/etc/sysConfig/UIResource/R2.8/test/Test11.png
/etc/sysConfig/UIResource/K1/uploadtips.png
/etc/sysConfig/UIResource/R2.8/checkbox_green.png
/etc/sysConfig/UIResource/R2.8/home_wifi.png
/etc/sysConfig/UIResource/R2.8/wifi_lock.png
/etc/sysConfig/UIResource/R2.8/test/Test19.png
/etc/sysConfig/UIResource/R2.8/dog.png
/etc/sysConfig/UIResource/K1/udiskicon.png
/etc/sysConfig/UIResource/R2.8/auto_refill_prev_page.png
/etc/sysConfig/UIResource/R2.8/EstateNormal.png
/etc/sysConfig/UIResource/R2.8/img_f011_refill_dir.png
/etc/sysConfig/UIResource/K1/wifi_state_1.png
/etc/sysConfig/UIResource/R2.8/img_f011_pause.png
/etc/sysConfig/UIResource/R2.8/test/Test20.png
/etc/sysConfig/defData/Auto_Em.gcode
/etc/sysConfig/UIResource/R2.8/cool_sub.png
/etc/sysConfig/UIResource/R2.8/kb_confirm_icon.png
/etc/sysConfig/UIResource/R2.8/img_f011_light_opened.png
/etc/sysConfig/UIResource/K1/silentmode.png
/etc/sysConfig/UIResource/R2.8/camera_enter.png
/etc/sysConfig/UIResource/R2.8/test/Test7.png
/etc/sysConfig/audio/chinese/1.sound_button_technology.wav
/etc/sysConfig/UIResource/R2.8/wifi_check.png
/etc/sysConfig/UIResource/R2.8/img_f011_humi_highest.png
/etc/sysConfig/UIResource/R2.8/img_f011_object_canceled.png
/etc/sysConfig/UIResource/K1/selected_state.png
/etc/sysConfig/defData/laser_offset_correction.gcode
/etc/sysConfig/UIResource/K1/lightopen.png
/etc/sysConfig/UIResource/R2.8/selftest_waiting.png
/etc/sysConfig/UIResource/R2.8/CFS_illustration_B.png
/etc/sysConfig/UIResource/K1/warningTips.png
/etc/sysConfig/UIResource/R2.8/prev.png
/etc/sysConfig/UIResource/K1/checking.gif
/etc/sysConfig/UIResource/R2.8/wifi_state_2.png
/etc/sysConfig/UIResource/R2.8/next_page_enable.png
/etc/sysConfig/UIResource/K1/K1_Machine.png
/etc/sysConfig/UIResource/R2.8/icon_rfid_flag.png
/etc/sysConfig/UIResource/R2.8/file_bed.png
/etc/sysConfig/UIResource/R2.8/file_next2.png
/etc/sysConfig/UIResource/R2.8/test/Test12.png
/etc/sysConfig/UIResource/K1/img_boot4.png
/etc/sysConfig/UIResource/K1/auto_refill_prev_page.png
/etc/sysConfig/UIResource/R2.8/img_f011_page_faq.png
/etc/sysConfig/UIResource/R2.8/setting_network.png
/etc/sysConfig/UIResource/R2.8/home_control.png
/etc/sysConfig/UIResource/K1/wifi_state_2.png
/etc/sysConfig/UIResource/R2.8/light_off.png
/etc/sysConfig/UIResource/R2.8/test/Test21.png
/etc/sysConfig/UIResource/K1/img_state_normal.png
/etc/sysConfig/UIResource/R2.8/img_f011_invalid_qrcode.png
/etc/sysConfig/audio/chinese/1.sound_button_short.wav
/etc/sysConfig/UIResource/R2.8/EstateIn.png
/etc/sysConfig/UIResource/R2.8/file.png
/etc/sysConfig/UIResource/R2.8/navigation_style.png
/etc/sysConfig/UIResource/K1/SortUpUnSelect.png
/etc/sysConfig/UIResource/K1/speed_mode.png
/etc/sysConfig/UIResource/R2.8/test/Test8.png
/etc/sysConfig/UIResource/R2.8/next_page_disable.png
/etc/sysConfig/UIResource/R2.8/auto_refill_next_page.png
/etc/sysConfig/UIResource/R2.8/setting_cloud.png
/etc/sysConfig/UIResource/K1/img_multi_color_guide_3.png
/usr/bin/display-server
/etc/sysConfig/defData/Auto_strain_gauge_pa_empty_print.gcode
/etc/sysConfig/UIResource/R2.8/invalid_texture_100x51.png
/etc/sysConfig/UIResource/R2.8/CFS_illustration_C.png
/etc/sysConfig/UIResource/R2.8/errorTips.png
/etc/sysConfig/UIResource/R2.8/img_f011_auto.png
/etc/sysConfig/UIResource/R2.8/file_prev2.png
/etc/sysConfig/UIResource/K1/filamentTip.png
/etc/sysConfig/UIResource/R2.8/cool_add.png
/etc/sysConfig/UIResource/R2.8/img_f011_title_tab_left.png
/etc/sysConfig/UIResource/R2.8/img_humi_bg.png
/etc/sysConfig/UIResource/K1/img_humidity_little_height_rangle.png
/etc/sysConfig/UIResource/R2.8/img_f011_page_filament.png
/etc/sysConfig/UIResource/R2.8/wifi_state_3.png
/etc/sysConfig/UIResource/K1/img_humidity_height_rangle.png
/etc/sysConfig/UIResource/R2.8/net_lock.png
/etc/sysConfig/UIResource/R2.8/root_info.png
/etc/sysConfig/UIResource/R2.8/test/Test13.png
/etc/sysConfig/UIResource/K1_Max/K2Plus_Checktips.png
/etc/sysConfig/UIResource/K1_Max/K1Max_Checktips.png
/etc/sysConfig/UIResource/K1/K1C_Machine.png
/etc/sysConfig/UIResource/R2.8/img_f011_humi_high.png
/usr/bin/master-server
/etc/sysConfig/UIResource/K1/wifi_state_3.png
/etc/sysConfig/UIResource/R2.8/test/Test22.png
/etc/sysConfig/defData/Auto_pressure_advance_multi.gcode
/etc/sysConfig/UIResource/R2.8/img_f011_nozzle.png
/etc/sysConfig/UIResource/K1_Max/img_edit.png
/etc/sysConfig/UIResource/K1/SortSelect.png
/etc/sysConfig/UIResource/K1/Add.png
/etc/sysConfig/UIResource/R2.8/history_cross.png
/etc/sysConfig/UIResource/R2.8/setting_guide_machine.png
/etc/sysConfig/UIResource/R2.8/test/Test1.png
/etc/sysConfig/UIResource/R2.8/file_hotend.png
/etc/sysConfig/UIResource/R2.8/setting_root.png
/etc/sysConfig/UIResource/K1/auto_refill_next_page.png
/etc/sysConfig/UIResource/R2.8/key_ok.png
/etc/sysConfig/UIResource/R2.8/Stepper_release.png
/etc/sysConfig/UIResource/R2.8/img_f011_object_focussed.png
/usr/bin/app-server
/etc/sysConfig/UIResource/R2.8/test/Test9.png
/etc/sysConfig/UIResource/R2.8/img_f011_next.png
/etc/sysConfig/UIResource/R2.8/EstateRunout.png
/etc/sysConfig/UIResource/R2.8/CFS_illustration_D.png
/etc/sysConfig/UIResource/K1/img_cfs_humidity_height.png
/etc/sysConfig/UIResource/R2.8/img_f011_page_file.png
/etc/sysConfig/UIResource/K1/zoffset_Add_Focus.png
/etc/sysConfig/UIResource/K1/print_completed.png
/etc/sysConfig/UIResource/K1/zoffset_Sub_Focus.png
/etc/sysConfig/UIResource/K1_Max/img_polygon.png
/etc/sysConfig/UIResource/R2.8/file_next.png
/etc/sysConfig/UIResource/R2.8/wifi_state_4.png
/etc/sysConfig/UIResource/R2.8/home_setting.png
/etc/sysConfig/UIResource/K1/wait.png
/etc/sysConfig/defData/Auto_pressure_advance_testpadvance.gcode
/etc/sysConfig/UIResource/R2.8/test/Test14.png
/etc/sysConfig/UIResource/R2.8/setting_brightness.png
/etc/sysConfig/UIResource/K1/stop_print_warning.png
/etc/sysConfig/UIResource/K1/img_cfs_humidity_little_height.png
/etc/sysConfig/UIResource/K1/SortUpSelect.png
/etc/sysConfig/UIResource/K1/remainingtimeicon.png
/etc/sysConfig/defData/laser_test_line.gcode
/etc/sysConfig/UIResource/K1/selfchecking_0.png
/etc/sysConfig/UIResource/R2.8/kb_delete_icon.png
/etc/sysConfig/defData/Auto_strain_gauge_pa_compute_config.json
/etc/sysConfig/UIResource/K1/wifi_state_4.png
/etc/sysConfig/UIResource/R2.8/light_on.png
/etc/sysConfig/UIResource/R2.8/test/Test23.png
/etc/sysConfig/UIResource/R2.8/print_opt.png
/etc/sysConfig/UIResource/R2.8/test/Test2.png
/etc/sysConfig/UIResource/R2.8/file_material.png
/etc/sysConfig/UIResource/K1/objectStateNormal.png
/etc/sysConfig/UIResource/R2.8/EstateOut.png
/etc/sysConfig/UIResource/R2.8/img_f011_bed.png
/etc/sysConfig/UIResource/R2.8/setting_language.png
/etc/sysConfig/UIResource/K1/gcode_default_small.png
/etc/sysConfig/UIResource/R2.8/img_f011_prev.png
/etc/sysConfig/defData/auto_laser_test.gcode
