# usr/lib/opkg/info/mcu_update.list
/usr/bin/mcu_util
/usr/bin/mcu_util_485
/etc/init.d/mcu_update
