# usr/lib/opkg/info/optee-client-3.7.list
/usr/lib/libteec.so.1.0
/usr/lib/libteec.so.1.0.0
/usr/sbin/tee-supplicant
/etc/init.d/tee-supplicant
/usr/lib/libteec.a
/usr/lib/libteec.so
/usr/lib/libteec.so.1
