# usr/lib/opkg/info/swupdate.list
/sbin/swupdate
/sbin/swupdate_cmd.sh
/etc/init.d/S99swupdate_autorun
/sbin/swupdate-progress
/sbin/swupdate-client
