# usr/lib/opkg/info/wipe_data.list
/etc/init.d/wipe_data
/usr/bin/wipe_data
